# spec/resources/stubs/initd/custom-app-bar
#!/bin/sh

# Autogenerated by foreman

### BEGIN INIT INFO
# Provides:             custom-app-bar
# Required-Start:       $syslog
# Required-Stop:        $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    custom-app-bar
# Description:          custom-app-bar
### END INIT INFO

NAME=custom-app-bar
DESC=custom-app-bar
DAEMON=/tmp/app/bar-script
DAEMON_USER=dummy
PIDDIR=/var/run/custom-app-bar
PIDFILE=/var/run/custom-app-bar/custom-app-bar.pid

test -x ${DAEMON} || exit 0
set -e
. /lib/lsb/init-functions


case "${1}" in
	start)
		log_daemon_msg "Starting ${DESC}" "${NAME}"
		mkdir -p ${PIDDIR}
		if (start-stop-daemon --start --make-pidfile --background --pidfile $PIDFILE --chuid $DAEMON_USER --exec $DAEMON -- argument1); then
			log_end_msg 0
		else
			log_end_msg 1
		fi
	;;
	stop)
		log_daemon_msg "Stopping ${DESC}" "${NAME}"
		if (start-stop-daemon --stop --oknodo --retry 20 --pidfile $PIDFILE --chuid $DAEMON_USER); then
			log_end_msg 0
		else
			log_end_msg 1
		fi
	;;
	status)
		status_of_proc -p ${PIDFILE} ${DAEMON} ${NAME}
	;;
	restart)
		${0} stop
		${0} start
	;;
	*)
		echo "Usage: ${0} {start|stop|status|restart}" >&2
	;;
esac

exit 0
